# crontab
0 * * * * /app/main >> /app/output/cron.log 2>&1